# deploy/command.service.d/10-immutable-release.conf
# ADR: the service must never read the working directory. It runs from an
# immutable release tree; `current` is repointed by scripts/release.sh via
# rename(2). A stray `next build` in the repo can no longer corrupt what is
# being served. (2026-07-12 split-brain outage.)
[Service]
WorkingDirectory=/opt/workspace/runtime/releases/command/current
